# linux/systemd/rc.local



```bash
touch "/etc/rc.local"
chmod +x "/etc/rc.local"

sudo systemctl "edi"t --full "rc-local"
```
